# wire-plug/wire-hdmi.service
[Unit]
Description=Wire Plug HDMI

[Service]
Type=oneshot
ExecStart=/usr/bin/wire-plug hdmi

[Install]
WantedBy=wire-hdmi.target
