# ctlbvs/90-backlight.rules
ACTION=="add", SUBSYSTEM=="backlight", RUN+="/bin/chmod 766 /sys/class/backlight/%k/brightness"